# deploy/systemd/energy-forecast-trigger.timer
[Unit]
Description=Dispatch energy forecasting daily workflow from Station at morning wake
Requires=energy-forecast-trigger.service

[Timer]
OnCalendar=*-*-* 06:55:00
Persistent=true
AccuracySec=1min
Unit=energy-forecast-trigger.service

[Install]
WantedBy=timers.target
