# deploy/systemd/energy-forecast-trigger.service
[Unit]
Description=Dispatch energy forecasting daily GitHub Actions workflow
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=smnfrs
Group=smnfrs
WorkingDirectory=/home/smnfrs/projects/energy-forecasting
Environment="HOME=/home/smnfrs"
Environment="XDG_CONFIG_HOME=/home/smnfrs/.config"
Environment="PATH=/home/smnfrs/.local/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
ExecStart=/home/smnfrs/projects/energy-forecasting/scripts/ops/trigger_daily_forecast.sh
StandardOutput=journal
StandardError=journal
TimeoutStartSec=180
